# .config/systemd/user/offlineimap-oneshot.service
[Unit]
Description=Get emails via IMAP

[Service]
Type=oneshot
WatchdogSec=300
ExecStart=offlineimap

[Install]
WantedBy=default.target