# .config/systemd/user/offlineimap-oneshot.timer
[Unit]
Description=Run offlineimap every 5min

[Timer]
OnUnitInactiveSec=5m

[Install]
WantedBy=timers.target